# llama-cpp.service
[Unit]
Description=llama.cpp server (profile-based)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=llama-cpp
Group=llama-cpp
ExecStart=/var/lib/llama-cpp/serve.sh
Restart=on-failure
RestartSec=5

# Uncomment if you have an NVIDIA GPU and need device access:
# SupplementaryGroups=video render

# Hardening (relax if it breaks GPU/driver access)
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=false
PrivateTmp=true

# Paths the service must access. /var/lib/llama-cpp holds the app + profiles.
# If models live on a separate disk and /var/lib/llama-cpp/models is a symlink
# to it (e.g. /data/var/lib/llama-cpp/models), list that target too so it stays
# accessible under stricter ProtectSystem settings. Drop the /data line if you
# keep everything on the root disk.
ReadWritePaths=/var/lib/llama-cpp
ReadWritePaths=-/data/var/lib/llama-cpp

[Install]
WantedBy=multi-user.target
